# templates/grafana.service.j2
[Unit]
Description=Grafana Container
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=on-failure
ExecStop=/usr/bin/docker stop {{ metrics_grafana_container_name }}
ExecStart=/usr/bin/docker start -a {{ metrics_grafana_container_name }}

[Install]
WantedBy=multi-user.target

